# services/systemd/olympus-order-listener.service
[Unit]
After=network.service

[Service]
User=cooper
ExecStart=/bin/bash -c '~/olympus/services/shell/olympus-order-listener.sh'

[Install]
WantedBy=multi-user.target